# supervisor/supervisord.conf
[supervisord]
[program:watch]
command=watch date
