# ansible/roles/projectify_app/templates/gunicorn.service.j2
# SPDX-License-Identifier: MIT
# SPDX-FileCopyrightText: 2009-2026 (c) Benoît Chesneau <benoitc@gunicorn.org>
# SPDX-FileCopyrightText: 2009-2015 (c) Paul J. Davis <paul.joseph.davis@gmail.com>

[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
Type=notify
NotifyAccess=main
User=projectify
Group=projectify
WorkingDirectory={{ app_root }}
Environment="SOCKET=/run/gunicorn.sock"
# TODO remove SITE_TITLE var
Environment="SITE_TITLE=Projectify Production on Hetzner"
Environment="DJANGO_SETTINGS_MODULE=projectify.settings.hetzner"
Environment="DJANGO_CONFIGURATION=Hetzner"
# TODO hard-code in Hetzner class
Environment="DATABASE_URL=postgres://%%2Frun%%2Fpostgresql/projectify"
Environment="STATIC_ROOT={{ static_root }}"
Environment="MEDIA_ROOT={{ media_root }}"
# TODO use LoadCredential for secrets
Environment="CREDENTIALS_FILE={{ credentials_file }}"
# TODO configure logs
ExecStart={{ venv_path }}/bin/python -m gunicorn
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# TODO add hardening
# Should not change files in /usr/share/projectify
# OK to access /var/lib/projectify/state/media
# Acceptable network stuff:
# - access unix socket (psql)
# - access/create unix socket (for caddy rev proxy)
# - make TCP requests to the outside world

[Install]
WantedBy=multi-user.target
